# deploy/backup-config.service
[Unit]
Description=Panel Configuration Backup Service
After=network.target

[Service]
Type=oneshot
User=panel
WorkingDirectory=/workspaces/panel
ExecStart=/usr/bin/python3 -c "
import sys
sys.path.insert(0, '/workspaces/panel')
from tasks import run_scheduled_config_backup
result = run_scheduled_config_backup()
print('Config backup result:', result)
"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target