# deploy/panel-logrotate.timer
[Unit]
Description=Daily logrotate for panel logs

[Timer]
OnCalendar=daily
Persistent=true
Unit=panel-logrotate.service

[Install]
WantedBy=timers.target
