# deploy/rq-worker.service
[Unit]
Description=RQ worker for panel (process background tasks)
After=network.target redis.service

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/home/YOUR_USER/panel
Environment="PATH=/home/YOUR_USER/panel/venv/bin"
ExecStart=/home/YOUR_USER/panel/venv/bin/python /home/YOUR_USER/panel/run_worker.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
